# server/browrts.service.example
[Unit]
Description=Ironfront RTS backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/BrowRts
Environment=HOST=127.0.0.1
Environment=PORT=2567
Environment=PUBLIC_WS_URL=wss://rts-api.example.com
ExecStart=/usr/bin/node /root/BrowRts/server/server.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
